# saic-java-mqtt-gateway/src/dist/systemd/saic-mqtt-gateway.service
[Unit]
Description=The SAIC MQTT Gateway
After=network.target
Wants=network-online.target

[Service]
Restart=always
Type=simple
#User=vmobil-dsp
#Group=vmobil-dsp
ExecStart=java -jar /opt/saic-mqtt-gateway/saic-mqtt-gateway.jar --config /etc/saic-mqtt-gateway.toml
#WorkingDirectory=/opt/advanova/vmobil-dsp
#Environment="JAVA_HOME=/usr/lib64/jvm/jre-1.8.0-openjdk"

[Install]
WantedBy=multi-user.target
